# deploy/com.local.qwen3tts.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.local.qwen3tts</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/aicer/qwen3tts/venv/bin/python</string>
        <string>server.py</string>
        <string>--host</string>
        <string>0.0.0.0</string>
        <string>--port</string>
        <string>9999</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/aicer/Library/Logs/qwen3tts.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/aicer/Library/Logs/qwen3tts.err</string>

    <key>WorkingDirectory</key>
    <string>/Users/aicer/qwen3tts/app</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/aicer/qwen3tts/venv/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>MLX_DEFAULT_DEVICE</key>
        <string>gpu</string>
        <key>LANG</key>
        <string>en_US.UTF-8</string>
        <key>QWEN3TTS_HOME</key>
        <string>/Users/aicer/qwen3tts</string>
        <key>QWEN3TTS_MODEL</key>
        <string>0.6b-custom</string>
    </dict>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
