# deploy/llkc.crontab
# LLM Knowledge Curator v2 incremental pipeline: daily at 06:00 local time.
0 6 * * * /bin/bash "/Users/aicer/Documents/Project/llm-knowledge-curator/scripts/cron_incremental_v2.sh"

# Qwen3TTS/STT health check: every 15 minutes, auto-restart if down.
*/15 * * * * /bin/bash "/Users/aicer/Documents/Project/llm-knowledge-curator/scripts/health_check.sh"

# Daily Brief: 07:00 generate summary of new items + distribution + actions
0 7 * * * /Library/Frameworks/Python.framework/Versions/3.13/bin/python3 /Users/aicer/Documents/Project/llm-knowledge-curator/scripts/cli.py brief >> /Users/aicer/Documents/Project/llm-knowledge-curator/output/cron_logs/brief.log 2>&1
